# .config/systemd/user/compton.service
[Unit]
Description=Compton

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=xorg.target
