# .config/systemd/user/nm-applet.service
[Unit]
Description=nm-applet

[Service]
ExecStart=/usr/bin/dbus-launch /usr/bin/nm-applet --sm-disable

[Install]
WantedBy=xorg.target
