# .config/systemd/user/setxkbmap.service
[Unit]
Description=setxkbmap

[Service]
Type=oneshot
ExecStart=/usr/bin/setxkbmap

[Install]
WantedBy=xorg.target
