# .config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper

[Service]
Type=oneshot
ExecStart=/bin/bash %h/.local/bin/wallpaper

[Install]
WantedBy=window-manager.target
